# docker/logpanel/supervisord.conf
[supervisord]
# nodaemon=true is REQUIRED: supervisord must run in the foreground as PID 1.
# With nodaemon=false it would fork and the container's main process would exit.
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ---------------------------------------------------------------------------
# Elasticsearch: non-root user, longer start window + more retries so Kibana
# (which has no real "depends") can come up after ES is ready.
# ---------------------------------------------------------------------------
[program:elasticsearch]
command=/opt/elasticsearch/bin/elasticsearch
user=elasticsearch
environment=ES_JAVA_OPTS="%(ENV_ES_JAVA_OPTS)s"
autostart=true
autorestart=true
startsecs=30
startretries=10
stopwaitsecs=30
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

# ---------------------------------------------------------------------------
# Kibana: non-root user, retries aggressively until ES is reachable.
# ---------------------------------------------------------------------------
[program:kibana]
command=/opt/kibana/bin/kibana
user=kibana
environment=ELASTICSEARCH_HOSTS="http://localhost:9200",SERVER_HOST="0.0.0.0",HOME="/opt/kibana",XDG_CONFIG_HOME="/opt/kibana/.config",XDG_CACHE_HOME="/opt/kibana/.cache"
autostart=true
autorestart=true
startsecs=15
startretries=20
stopwaitsecs=30
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

# ---------------------------------------------------------------------------
# nginx: auth gate (basic-auth + rate limit + TLS). Runs in foreground.
# ---------------------------------------------------------------------------
[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
startsecs=15
startretries=5
stopwaitsecs=10
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
